# centreon-web/files/supervisord/centreon-core.ini
[program:centcore]
priority=1
user=centreon
command=/usr/bin/perl /usr/share/centreon/bin/centcore --logfile=/var/log/centreon/centcore.log --severity=error --config=/etc/centreon/conf.pm
autostart=true
stdout_logfile=/var/log/%(program_name)s-sc-out.log
stdout_logfile_maxbytes=64MB
stdout_logfile_backups=7
stdout_events_enabled=false
stderr_logfile=/var/log/%(program_name)s-sc-err.log
stderr_logfile_maxbytes=64MB
stderr_logfile_backups=7
stderr_events_enabled=false
