# centreon-web/files/supervisord/php-fpm.ini
[program:php7-fpm]
priority=1
command=/opt/rh/rh-php71/root/usr/sbin/php-fpm -F 
autostart=true
stdout_logfile=/var/log/%(program_name)s-sc-out.log
stdout_logfile_maxbytes=64MB
stdout_logfile_backups=7
stdout_events_enabled=false
stderr_logfile=/var/log/%(program_name)s-sc-err.log
stderr_logfile_maxbytes=64MB
stderr_logfile_backups=7
stderr_events_enabled=false
