# centreon-engine/files/supervisord/postfix.ini
[program:postfix]
process_name=postfix
command=/usr/local/sbin/postfix-wrapper.sh
startsecs=0
autorestart=false
